# .config/systemd/user/clean_journal.service
[Unit]
Description=Clean journal directory

[Service]
Type=simple
ExecStart=/home/ed/bin/clean_journal.sh
