# .config/systemd/user/del_cache.service
[Unit]
Description=Delete ~/.cache/* 

[Service]
Type=simple
ExecStart=/home/ed/bin/del_cache.sh
