# .config/systemd/user/mirrorlist.timer
[Unit]
Description=Mirrorlist update weekly timer

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
