# .config/systemd/user/orphan_pack.timer
[Unit]
Description=Clean orphan packages monthly

[Timer]
OnCalendar=monthly
#OnCalendar=*-*-* 08:45:00
Persistent=true

[Install]
WantedBy=timers.target
