# .config/systemd/user/pac_cache.service
[Unit]
Description=Clean package cache

[Service]
Type=simple
ExecStart=/home/ed/bin/pac_cache.sh
